# crontab
0 0 * * * ~/download_copy_covid_colombia.sh
